# port/pkg-plist
@mode 04755
bin/crate
@mode
man/man5/crate.5.gz
%%DAEMON%%sbin/crated
%%DAEMON%%@mode 0555
%%DAEMON%%etc/rc.d/crated
%%DAEMON%%@mode
%%DAEMON%%@sample etc/crated.conf.sample
%%SNMPD%%sbin/crate-snmpd
%%SNMPD%%share/snmp/mibs/CRATE-MIB.txt
%%COMPLETIONS%%share/crate/completions/crate.sh
%%COMPLETIONS%%@dir share/crate/completions
%%COMPLETIONS%%@dir share/crate
%%EXAMPLES%%share/examples/crate/amass.yml
%%EXAMPLES%%share/examples/crate/aria2.yml
%%EXAMPLES%%share/examples/crate/blender-gpu.yml
%%EXAMPLES%%share/examples/crate/broken/kodi.yml
%%EXAMPLES%%share/examples/crate/broken/midori.yml
%%EXAMPLES%%share/examples/crate/broken/opera.yml
%%EXAMPLES%%share/examples/crate/broken/qbittorrent+tor.yml
%%EXAMPLES%%share/examples/crate/broken/zeronet+tor.yml
%%EXAMPLES%%share/examples/crate/broken/zeronet.yml
%%EXAMPLES%%share/examples/crate/chromium+i2pd.yml
%%EXAMPLES%%share/examples/crate/chromium+tor.yml
%%EXAMPLES%%share/examples/crate/chromium-gpu.yml
%%EXAMPLES%%share/examples/crate/chromium-headless.yml
%%EXAMPLES%%share/examples/crate/chromium-privacy.yml
%%EXAMPLES%%share/examples/crate/chromium.yml
%%EXAMPLES%%share/examples/crate/fetch.yml
%%EXAMPLES%%share/examples/crate/firefox+i2pd.yml
%%EXAMPLES%%share/examples/crate/firefox-gpu.yml
%%EXAMPLES%%share/examples/crate/firefox-headless.yml
%%EXAMPLES%%share/examples/crate/firefox.yml
%%EXAMPLES%%share/examples/crate/gimp-headless.yml
%%EXAMPLES%%share/examples/crate/gimp.yml
%%EXAMPLES%%share/examples/crate/glxgears-gpu.yml
%%EXAMPLES%%share/examples/crate/glxgears.yml
%%EXAMPLES%%share/examples/crate/gogs.yml
%%EXAMPLES%%share/examples/crate/gzip.yml
%%EXAMPLES%%share/examples/crate/i2pd.yml
%%EXAMPLES%%share/examples/crate/libreoffice-headless.yml
%%EXAMPLES%%share/examples/crate/matrix-coturn.yml
%%EXAMPLES%%share/examples/crate/matrix-nginx.yml
%%EXAMPLES%%share/examples/crate/matrix-synapse.yml
%%EXAMPLES%%share/examples/crate/matrix/README.md
%%EXAMPLES%%share/examples/crate/matrix/matrix-backup.sh
%%EXAMPLES%%share/examples/crate/matrix/matrix-health.sh
%%EXAMPLES%%share/examples/crate/meld.yml
%%EXAMPLES%%share/examples/crate/mpv.yml
%%EXAMPLES%%share/examples/crate/multi-interface.yml
%%EXAMPLES%%share/examples/crate/nginx-netgraph.yml
%%EXAMPLES%%share/examples/crate/nginx-passthrough.yml
%%EXAMPLES%%share/examples/crate/nginx-vlan.yml
%%EXAMPLES%%share/examples/crate/nginx.yml
%%EXAMPLES%%share/examples/crate/nmap.yml
%%EXAMPLES%%share/examples/crate/postgres.yml
%%EXAMPLES%%share/examples/crate/qbittorrent+i2pd.yml
%%EXAMPLES%%share/examples/crate/qbittorrent.yml
%%EXAMPLES%%share/examples/crate/qtox.yml
%%EXAMPLES%%share/examples/crate/redis.yml
%%EXAMPLES%%share/examples/crate/stack-matrix.yml
%%EXAMPLES%%share/examples/crate/stack-webapp.yml
%%EXAMPLES%%share/examples/crate/syncthing.yml
%%EXAMPLES%%share/examples/crate/telegram-desktop.yml
%%EXAMPLES%%share/examples/crate/thunderbird.yml
%%EXAMPLES%%share/examples/crate/tor.yml
%%EXAMPLES%%share/examples/crate/vlc.yml
%%EXAMPLES%%share/examples/crate/webapp.yml
%%EXAMPLES%%share/examples/crate/wget.yml
%%EXAMPLES%%share/examples/crate/xeyes.yml
%%EXAMPLES%%share/examples/crate/xfce-desktop.yml
%%EXAMPLES%%share/examples/crate/yt-dlp.yml
%%EXAMPLES%%@dir share/examples/crate/broken
%%EXAMPLES%%@dir share/examples/crate/matrix
%%EXAMPLES%%@dir share/examples/crate
